# monitors/99-monitors-hotplug.rules
ACTION=="change" \
, SUBSYSTEM=="drm"\
, RUN+="/bin/systemctl restart monitors.service"
